# arena-srt-relay.service
[Unit]
Description=Arena SRT audio relay (PCM->AAC, port 8890->8892)
After=network.target mediamtx.service

[Service]
Type=simple
ExecStart=/opt/arena/srt-audio-relay.sh
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
